# deploy/cityquest.service
[Unit]
Description=CityQuest game server
After=network.target

[Service]
Type=simple
User=cityquest
WorkingDirectory=/opt/cityquest
ExecStart=/opt/cityquest/cityquest
Restart=on-failure
RestartSec=5

Environment=DB_PATH=/opt/cityquest/data/local.db
Environment=SPA_DIR=/opt/cityquest/web
Environment=HTTP_ADDR=:8443
Environment=TLS_CERT=/opt/cityquest/tls/cert.pem
Environment=TLS_KEY=/opt/cityquest/tls/key.pem

# Hardening
ProtectSystem=strict
ReadWritePaths=/opt/cityquest/data
ReadOnlyPaths=/opt/cityquest/tls
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
